# bin/jiawei.tomato.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>jiawei.github.tomato</string>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jiawei/workspace/github/tomato/bin/run.sh</string>
    </array>
    <key>UserName</key>
    <string>jiawei</string>
    <key>GroupName</key>
    <string>staff</string>
</dict>
</plist>
